# 0logs/dropbear-2014.66-i686_slk600.files
dropbear-install/etc/init.d/rc.dropbear
dropbear-install/usr/bin/dbclient
dropbear-install/usr/bin/dropbearconvert
dropbear-install/usr/bin/dropbearkey
dropbear-install/usr/bin/dropbearmulti
dropbear-install/usr/bin/scp
dropbear-install/usr/bin/ssh
dropbear-install/usr/sbin/dropbear
dropbear-install/usr/share/man/man1/dbclient.1
dropbear-install/usr/share/man/man1/dropbearconvert.1
dropbear-install/usr/share/man/man1/dropbearkey.1
dropbear-install/usr/share/man/man8/dropbear.8
